# infra/ansible/templates/speako.service.j2
[Unit]
Description=Comments plugin

[Service]
User=root
WorkingDirectory=/opt/
ExecStart=/opt/aawaz
Restart=always

[Install]
WantedBy=multi-user.target
